# flighttrak-dashboard.service
[Unit]
Description=FlightTrak Enhanced Dashboard
After=network.target flightalert.service
Wants=flightalert.service

[Service]
Type=simple
User=kurt
Group=kurt
WorkingDirectory=/home/kurt/flighttrak
Environment="PATH=/home/kurt/flighttrak/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/kurt/flighttrak/venv/bin/python3 /home/kurt/flighttrak/enhanced_dashboard.py
Restart=always
RestartSec=10
StandardOutput=append:/home/kurt/flighttrak/dashboard.log
StandardError=append:/home/kurt/flighttrak/dashboard.log

[Install]
WantedBy=multi-user.target