# flighttrak-ai-intelligence.service
[Unit]
Description=FlightTrak AI Intelligence Service - Advanced Event Detection
After=network.target
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=kurt
Group=kurt
WorkingDirectory=/home/kurt/flighttrak
Environment="PATH=/home/kurt/flighttrak/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/home/kurt/flighttrak"

# Main AI Intelligence Service
ExecStart=/home/kurt/flighttrak/venv/bin/python3 /home/kurt/flighttrak/ai_intelligence_service.py

# Service management
ExecStartPre=/bin/sleep 5
ExecReload=/home/kurt/flighttrak/venv/bin/python3 /home/kurt/flighttrak/send_service_notification.py reload
ExecStopPost=/home/kurt/flighttrak/venv/bin/python3 -c "import logging; logging.basicConfig(filename='/home/kurt/flighttrak/ai_intelligence_service.log', level=logging.INFO); logging.info('AI Intelligence Service stopped')"

# Restart configuration
Restart=always
RestartSec=30
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=append:/home/kurt/flighttrak/ai_intelligence_service.log
StandardError=append:/home/kurt/flighttrak/ai_intelligence_service.log

# Resource limits
MemoryLimit=1G
CPUQuota=200%

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/kurt/flighttrak

[Install]
WantedBy=multi-user.target